# scripts/network-monitor.service
[Unit]
Description=Network Monitor TUI Service
After=network.target
Wants=network.target
Documentation=man:network-monitor(1)

[Service]
Type=simple
ExecStart=/usr/bin/network-monitor -c /etc/network-monitor.conf
Restart=always
RestartSec=5
User=root
Group=root

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true

# Resource limits
MemoryMax=512M
CPUQuota=80%

# Capabilities required for packet capture
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
